# Norns/packaging/sidekick_pkg/postinst
#if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ]; then
#fi
udevadm control --reload-rules >/dev/null || true



USERDIR=/home/we
SKDIR="$USERDIR/sidekick"
su we -c "mkdir -p $SKDIR $SKDIR/patches $SKDIR/system"
su we -c "cp -r /usr/local/sidekick/patches $SKDIR"
su we -c "cp -r /usr/local/sidekick/pdexternals $SKDIR"
su we -c "cp -r /usr/local/sidekick/*.json $SKDIR"
su we -c "cp -r /usr/local/sidekick/*.sh $SKDIR"

su we -c "mkdir -p $SKDIR/.config/SuperCollider/"
su we -c "cp /usr/local/sidekick/startup.scd  $SKDIR/.config/SuperCollider/"

# stop norns 
systemctl stop norns.target
systemctl disable norns.target

systemctl enable sidekick.target
systemctl restart sidekick.target
systemctl restart sidekick
