# Bela/packaging/mec_pkg/postinst
#if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ]; then
#fi
systemctl enable mec
systemctl start mec

rm /root/MEC
ln -s /usr/local/MEC /root/
