# patchbox/packaging/mec_pkg/postinst
#if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ]; then
#fi
udevadm control --reload-rules >/dev/null || true

#start with pisound btn
cp /usr/local/MEC/pisound.conf /etc
#systemctl enable mec
#systemctl start mec
