# TT/packaging/sidekick_pkg/prerm
systemctl stop sidekick.target
systemctl disable sidekick.target


if [ $1 = "remove" ]; then
    USERDIR=/home/we
    SKDIR="$USERDIR/sidekick"
    rm -rf $SKDIR
    #systemctl enable norns.target
    #systemctl start norns.target
fi

