# TT/packaging/sidekick_pkg/prerm
systemctl stop sidekick.target
systemctl disable sidekick.target


USERDIR=/home/we
SKDIR="$USERDIR/sidekick"
rm -rf $SKDIR

