# Nebulae/packaging/mec_pkg/postinst
#if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ]; then
#fi
udevadm control --reload-rules >/dev/null || true


USERDIR=/home/alarm
cp /usr/local/MEC/mec.instr $HOME/instr

systemctl enable mec
#systemctl start mec
