# Nebulae/packaging/mec_pkg/prerm
systemctl stop mec
systemctl disable mec

USERDIR=/home/alarm
rm $USERDIR/instr/mec.instr