# PI/packaging/orac_pkg/prerm
systemctl stop orac
systemctl disable orac
