# PI/packaging/mec_pkg/prerm
systemctl stop mec
systemctl disable mec
